# crontab
*/15 * * * * /home/pi/TemperatureTrackerService/start.sh